# .crontab
*    * * * * $HOME/.dotfiles/misc/mail-sync
0    * * * * find ~ -name '*.core' -type f -delete
0    0 * * * find ~ -name '*~'     -type f -delete
0    0 * * * find ~ -name '#*'     -type f -delete
0    * * * * DISPLAY=:0 dunstify "$(date)"

